# com.shashank.switchbar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.shashank.switchbar</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/shashank.kk/.switchdeck-venv/bin/python3</string>
    <string>/Users/shashank.kk/Documents/GitHub/switchdeck/switchbar.py</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/shashank.kk/Documents/GitHub/switchdeck</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><false/>
  <key>StandardOutPath</key><string>/Users/shashank.kk/switchdeck.log</string>
  <key>StandardErrorPath</key><string>/Users/shashank.kk/switchdeck.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/shashank.kk/.local/bin</string>
  </dict>
</dict>
</plist>
